# crontab
* * * * * ubuntu cd /home/ubuntu/apps/oscars/repository && source ../virtualenv/bin/activate && ./build_awards_json.py
