# crontab
* * * * * ubuntu cd /home/ubuntu/apps/oscars/repository && source ../virtualenv/bin/activate && ../virtualenv/bin/fab deploy_awards_json
