# crontab
* * * * * ubuntu bash /home/ubuntu/apps/oscars/repository/cron.sh fab $DEPLOYMENT_TARGET deploy_awards_json
